# fs/etc/local.d/netm.start
#!/bin/sh

exec /usr/local/bin/netu
  # Можно сделать менее словесный формат вывода,
  # такой же как у других скриптов инициализации.
